# files/default/exhibitor.conf
# exhibitor - regular background program processing daemon
#
# exhibitor is a co-process for instance monitoring,
# backup/recovery, cleanup and visualization for zookeeper

description "zookeeper co-process"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/exhibitor/exhibitor.sh
